# supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
logfile_maxbytes=50MB
logfile_backups=10
loglevel=info
pidfile=/var/run/supervisord.pid

[program:inject-env]
command=/app/inject-env.sh
directory=/app
user=root
autostart=true
autorestart=false
startsecs=0
priority=100
stdout_logfile=/var/log/supervisor/inject-env.log
redirect_stderr=true

[program:backend]
command=/app/backend/main
directory=/app/backend
user=appuser
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=DB_PATH="/app/data/blog.db",GIN_MODE="debug",PORT="8085",APP_VERSION="%(ENV_APP_VERSION)s",BUILD_NUMBER="%(ENV_BUILD_NUMBER)s",BUILD_DATE="%(ENV_BUILD_DATE)s",GIT_COMMIT="%(ENV_GIT_COMMIT)s",GIT_BRANCH="%(ENV_GIT_BRANCH)s"
priority=200

[program:frontend]
command=node server.js
directory=/app/frontend
user=appuser
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=NODE_ENV="production",PORT="3000",HOSTNAME="0.0.0.0"
priority=300

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
user=root
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=400